# etc/systemd/system/modem-watchdog.service
[Unit]
Description=Modem watchdog - restart modem if not registered every 10min
After=network.target

[Service]
Type=simple
ExecStart=/opt/captive-portal/modem-watchdog.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
