# etc/systemd/system/provision-watchdog.service
[Unit]
Description=Fallback to AP if STA fails
After=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/opt/captive-portal/provision-watchdog.sh

[Install]
WantedBy=multi-user.target
